# scheduled-testing.service
[Unit]
Description=scheduled network testing
After=network.target

[Service]
User=nanopi
ExecStart=/usr/bin/python3 -u /home/nanopi/schedule.py
RestartSec=3
Restart=always

[Install]
WantedBy=multi-user.target
