# roles/testing/templates/sockperf-server.service
[Unit]
Description=sockperf server
After=network.target

[Service]
ExecStart=/home/adam/testing-server/sockperf server -p {{ sockperf_port }}
RestartSec=3
Restart=always

[Install]
WantedBy=multi-user.target
